# var/lib/dpkg/info/bb-customizations.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/bb-customizations
/usr/share/doc/bb-customizations/copyright
/usr/share/doc/bb-customizations/changelog.Debian.gz
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/dtbo
/usr/share/initramfs-tools/hooks/ti_pru_firmware
/etc
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/83-eqep-noroot.rules
/etc/udev/rules.d/82-gpio-config-pin.rules
/etc/udev/rules.d/60-omap-tty.rules
/etc/udev/rules.d/81-pwm-noroot.rules
/etc/udev/rules.d/80-gpio-noroot.rules
/etc/kernel
/etc/kernel/postinst.d
/etc/kernel/postinst.d/zz-uenv_txt
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/generic-board-startup.service
